# dot_config/systemd/user/backup.service.tmpl
[Unit]
Description=Backup multiple directories to target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/backup.sh

[Install]
WantedBy=default.target
