# dot_config/systemd/user/backup.timer
[Unit]
Description=Run directory sync every minute

[Timer]
OnCalendar=*:0/1
Persistent=true
Unit=backup.service

[Install]
WantedBy=timers.target
